# SOURCES/mapstory.supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock 

[supervisord]
pidfile=/var/run/supervisor.pid 
logfile_backups=1

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:mapstory]
command =  /var/lib/mapstory/.venv/bin/uwsgi --http 0.0.0.0:8000 
           --wsgi-file /var/lib/mapstory/mapstory-geonode/mapstory/wsgi.py 
           --processes 4 
           -t 120 
           --disable-logging -M 
           --need-app 
           -b 32768 
           --pp /var/lib/mapstory/.venv 
           --uid mapstory
           --max-request 5000
           --static-map /robots.txt=/var/lib/mapstory/mapstory-geonode/mapstory/templates/robots.txt
user=mapstory
stdout_logfile=/var/log/mapstory/stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/mapstory/stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5 
stopsignal=INT
